# deploy/supervisor/huey.conf
[program:huey]
command=/usr/local/bin/huey_consumer.py workers.huey -k process -w 4
directory=/SkyRadar
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true
