# deploy/supervisor/skyradar.conf
[program:skyradar]
command=/usr/local/bin/gunicorn -w10 -k uvicorn.workers.UvicornWorker -b0.0.0.0:8888 api:app
directory=/SkyRadar
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true
